# ppt18w/ppt18W.timer
[Unit]
Description=Run ppt18W.service every 2 seconds

[Timer]
OnBootSec=2s
OnUnitActiveSec=2s

[Install]
WantedBy=timers.target
